# deploy/systemd/kafil-backup.service
[Unit]
Description=Back up Kafil PostgreSQL and persistent storage off-VPS
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=kafil-deploy
Group=kafil-deploy
ExecStart=/opt/kafil/current/scripts/backupVps.sh
PrivateTmp=true
NoNewPrivileges=true
